# liewa-gui/systemd/liewa.service
[Unit]
Description=Liewa Service
[Service]
Type=simple
ExecStart=/home/lennart/.local/bin/liewa
