# liewa-gui/systemd/liewa.service
[Unit]
Description=Very useful script
[Service]
Type=simple
ExecStart=/home/lennart/.local/bin/liewa
